# .config/autostart/fix-res.desktop
[Desktop Entry]
Type=Application
Name=Fix Resolution
Exec=/home/aginiaux/.local/bin/set-resolution.sh
X-GNOME-Autostart-enabled=true
