# .config/autostart/rebinds.desktop
[Desktop Entry]
Type=Application
Name=Rebinds
Exec=/home/aginiaux/.local/bin/rebinds.sh
Hidden=false
X-GNOME-Autostart-enabled=true
